# usr/local/relianoid/share/relianoid.cron
# WARNING: Do not edit this file, it will be overwritten.

*/5  * * * *    root    /usr/local/relianoid/bin/noid-collector &>/dev/null
*/10 * * * *    root    /usr/local/relianoid/bin/noid-ntp &>/dev/null
*/30 * * * *    root    /usr/bin/find /tmp/ -name "cgisess_*" -mmin +60 -exec rm -f {} \;  &>/dev/null
*/30 * * * *    root    /usr/bin/find /tmp/ -name "support*.tar.gz" -mmin +60 -exec rm -f {} \;  &>/dev/null
0   22 * * *	root    /usr/local/relianoid/bin/noid-updater &>/dev/null
